# skg_deploy/scripts/skg-train.service
[Unit]
Description=SKG Daily Model Fine-Tune
After=network.target skg.service
Wants=skg.service

[Service]
Type=oneshot
User=%i
Environment=SKG_HOME=/opt/skg
Environment=SKG_STATE_DIR=/var/lib/skg
Environment=SKG_CONFIG_DIR=/etc/skg
Environment=PATH=/opt/skg/.venv/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/opt/skg
ExecStart=/opt/skg/.venv/bin/python -m skg.training.scheduler_main
TimeoutStartSec=18000
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skg-train

[Install]
WantedBy=default.target
